# files/etc/supervisor/conf.d/cron.conf.template
[program:cron]
command = cron -f
autostart = false
autrestart = true
